# embedded/systemd/camera-streamer.service
[Unit]
Description=Camera RTSP streaming service
After=network.target

[Service]
ExecStart=/usr/bin/camera-streamer
# Re-read the config and re-serve (e.g. switch listen= between usb/ethernet).
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
